# mission-control.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.mission-control</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>exec /Users/mattbruce/.openclaw/workspace/projects/Web/alex-mission-control/start-mission-control.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mattbruce/.openclaw/workspace/projects/Web/alex-mission-control</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/mattbruce/.openclaw/workspace/projects/Web/alex-mission-control/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mattbruce/.openclaw/workspace/projects/Web/alex-mission-control/logs/launchd-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/mattbruce</string>
    </dict>
</dict>
</plist>
